# src/lia/_debian-arm64/install.sh
#!/bin/bash

################################################################################
# Wireless
################################################################################
[ ! -f /etc/modprobe.d/blacklist-brcmfmac.conf ] && echo "blacklist brcmfmac" | tee -a /etc/modprobe.d/blacklist-brcmfmac.conf

################################################################################
# Samba
################################################################################
service smbd stop
systemctl disable smbd
service nmbd stop
systemctl disable nmbd
